# pi-gen_stage/stage2/04-install-packages/files/WebDisplay.service
[Unit]
Description=WebDisplay Service

[Install]
WantedBy=default.target

[Service]
Restart=always
ExecStart=python /WebDisplay/main.py