# pi-gen_stage/stage2/04-install-requirements/files/xorg.service
[Unit]
Description=Start X at boot for kiosk
After=multi-user.target

[Service]
User=pi
WorkingDirectory=/home/pi
Environment=DISPLAY=:0
ExecStart=/usr/bin/startx /home/pi/.xinitrc
Restart=always
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target
